# templates/device.path
[Path]
PathExists={{ item.pci_device }}
